# deploy/training-system.service
[Unit]
Description=Training System Web Dashboard v3.0
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/training-system
Environment=TRAIN_WEB_HOST=0.0.0.0
Environment=TRAIN_WEB_PORT=8082
Environment=ROOT_PATH=/training
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m uvicorn training.web.app:app --host 0.0.0.0 --port 8082
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
